# .docker/conf.d/supervisord.conf

[program:laravel-logs]
command=tail -qF /var/www/storage/logs/laravel.log
autostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
