# .docker/supervisor.d/queue-worker.conf
[program:queue-worker]
process_name=%(process_num)02d
command=/bin/supervisor-prefix-logs php /var/www/artisan queue:work --sleep=3 --tries=3 --max-time=86400
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=4
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
